# scripts/novabot.service
# user service
[Unit]
Description=NovaBot Service
Documentation=https://github.com/NovaBotDevs/NovaBot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/.local/share/novabot
ExecStart=/usr/bin/sh -c '/usr/bin/novabot run || { /usr/bin/novabot init && /usr/bin/novabot run; }'
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=novabot-%u
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
